# rootfs_making/seele-udev-seed.service
[Unit]
Description=Seed static udev database entries for Seele input devices
After=systemd-udevd.service systemd-udev-trigger.service
Before=multi-user.target graphical.target

[Service]
Type=oneshot
ExecStart=/usr/bin/seele-udev-seed.sh

[Install]
WantedBy=multi-user.target
